# system/services/emmi.target
[Unit]
Description=EMMI Facility

[Install]
WantedBy=multi-user.target
